# noobmon.cron
00,10,20,30,40,50 * * * * root /opt/noobmon/noobmon.sh
